# rc.local
#!/bin/sh -e
sudo /etc/ec2-prep.sh bridge
exit 0
