# tor-prep.conf
# installs and configures Tor on Ubuntu Amazon AMI
description     "Tor Cloud init script"

start on startup

task
exec /etc/ioerror.sh
